# etc/supervisor.d/celery.ini
[program:ingest]
command=/home/liorm/work/.virtualenvs/cloudify-cosmo/bin/celery worker -A celery_demo.celery_config -l debug -Q ingest -n ingest --purge
directory=/home/liorm/work/lior_repos/celery_demo
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/celery_ingest.log
logfile_maxbytes=10M
logfile_backups=5
startsecs=10
user=liorm
environment=CELERY_TASK=ingest

[program:operation]
command=/home/liorm/work/.virtualenvs/cloudify-cosmo/bin/celery worker -A celery_demo.celery_config -l debug -Q optimize -n optimize --purge
directory=/home/liorm/work/lior_repos/celery_demo
autostart=true
autorestart=true
redirect_stderr=true
logfile_maxbytes=10M
logfile_backups=5
startsecs=10
stdout_logfile=/tmp/celery_optimize.log
user=liorm
environment=CELERY_TASK=operation

[program:generic]
command=/home/liorm/work/.virtualenvs/cloudify-cosmo/bin/celery worker -A celery_demo.celery_config -l debug -Q default -n generic --purge
directory=/home/liorm/work/lior_repos/celery_demo
autostart=true
autorestart=true
redirect_stderr=true
logfile_maxbytes=10M
logfile_backups=5
startsecs=10
stdout_logfile=/tmp/celery.log
user=liorm
environment=CELERY_TASK=default


[group:celery]
programs=ingest, operation, generic
priority=2
