# etc/supervisor.d/celery.ini
[program:ingest]
command=/home/lior/work/.virtualenvs/celery_demo/bin/celery worker -A celery_demo.tasks -l debug -Q ingest -n ingest --purge
directory=/home/lior/maze_game/celery_demo
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/celery_ingest.log
logfile_maxbytes=10M
logfile_backups=5
startsecs=10
user=lior
environment=CELERY_TASK=ingest

[program:operation]
command=/home/lior/work/.virtualenvs/celery_demo/bin/celery worker -A celery_demo.tasks -l debug -Q operation -n operation --purge
directory=/home/lior/maze_game/celery_demo
autostart=true
autorestart=true
redirect_stderr=true
logfile_maxbytes=10M
logfile_backups=5
startsecs=10
stdout_logfile=/tmp/celery_operation.log
user=lior
environment=CELERY_TASK=operation

[program:generic]
command=/home/lior/work/.virtualenvs/celery_demo/bin/celery worker -A celery_demo.tasks -l debug -Q default -n default --purge
directory=/home/lior/maze_game/celery_demo
autostart=true
autorestart=true
redirect_stderr=true
logfile_maxbytes=10M
logfile_backups=5
startsecs=10
stdout_logfile=/tmp/celery.log
user=lior
environment=CELERY_TASK=default


[group:celery]
programs=ingest, operation, generic
priority=2
