# etc/supervisor.d/flower.ini
[program:flower]
command=/home/liorm/work/.virtualenvs/cloudify-cosmo/bin/celery -A celery_demo.celery_config flower --address=0.0.0.0
directory=/home/liorm/work/lior_repos
autostart=false
redirect_stderr=true
stdout_logfile=/tmp/flower.log
priority=3
logfile_maxbytes=10M
logfile_backups=5
user=liorm
startsecs=10
