# etc/supervisor.d/flower.ini
[program:flower]
command=/home/lior/work/.virtualenvs/celery_demo/bin/celery -A celery_demo.tasks flower --address=0.0.0.0 --broker_api=redis://localhost:26379/0 
directory=/home/lior/maze_game/celery_demo
autostart=true
redirect_stderr=true
stdout_logfile=/tmp/flower.log
priority=3
logfile_maxbytes=10M
logfile_backups=5
user=lior
startsecs=10
