# etc/supervisor.d/redis.ini
[program:redis]
command=/home/liorm/work/tools/redis-3.0.7/src/redis-server /home/liorm/work/etc/redis.conf
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/redis.log
priority=1
logfile_maxbytes=10M
user=liorm
logfile_backups=5
